# etc/supervisord.d/kafka-manager.ini
[program:kafka-manager]
command=/opt/qnib/kafka-manager/bin/start_kafka-manager.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
